# .config/systemd/user/upgradables.service
[Unit]
Description=Check upgradable packages and save

[Service]
Type=oneshot
ExecStart=%h/.config/awesome/scripts/check-upgradables.sh
