# .config/systemd/user/wallpaper.timer
[Unit]
Description=Timer for a wallpaper change

[Timer]
OnStartupSec=5min
OnUnitActiveSec=4min

[Install]
WantedBy=graphical-session.target
